# supervisor/conf.d/confd.conf
[program:confd]
command=/usr/local/bin/python main.py
directory=/root/confd
autostart=true
autorestart=true
redirect_stderr=true
priority=100
