# supervisor/conf.d/confd.conf
[program:confd]
command=/usr/local/bin/python /root/confd/main.py
autostart=false
autorestart=true
startsecs=5
redirect_stderr=true
