# supervisor/conf.d/haproxy.conf
[program:haproxy]
command=/usr/local/sbin/haproxy -f /usr/local/etc/haproxy/haproxy.cfg -n 4096 -N 4096
autostart=true
autorestart=true
user=root
